# 99-omen-rgb-keyboard.rules
# HP OMEN RGB Keyboard udev rules
# This allows users in the 'input' group to control RGB keyboard settings
# without requiring sudo privileges

# When the omen-rgb-keyboard platform device appears, set group and permissions
# on all sysfs attributes so non-root users in 'input' can write to them.
SUBSYSTEM=="platform", KERNEL=="omen-rgb-keyboard", ACTION=="add", RUN+="/bin/sh -c 'chgrp input /sys/devices/platform/omen-rgb-keyboard/rgb_zones/* && chmod 0664 /sys/devices/platform/omen-rgb-keyboard/rgb_zones/*'"
